# AWS_Intro_to_Ansible/Intermediate/Ansible/config_files/flaskapp.service
[Unit]
Description=Gunicorn daemon to serve my flaskapp
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/home/ec2-user/flask_app
ExecStart=/home/ec2-user/flask_app --bind 0.0.0.0:8000 app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target